# linux/systemctl/flask-app.service
[Unit]
Description=Flask Application Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/flask-app
Environment="PATH=/opt/flask-app/venv/bin"
ExecStart=/opt/flask-app/venv/bin/python /opt/flask-app/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target




